# document-embedding-api.service
[Unit]
Description=Document Embedding REST API Service with Integrated OpenAPI Documentation
After=network.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=chris
Group=chris
WorkingDirectory=/home/chris/Projects/EmbeddingsService
Environment=PATH=/home/chris/Projects/EmbeddingsService/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/chris/Projects/EmbeddingsService/venv/bin/python production_rest_api_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=document-embedding-api

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/chris/Projects/EmbeddingsService

[Install]
WantedBy=multi-user.target
